# deploy/autosub-worker.service
[Unit]
Description=AutoSub Celery worker
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=autosub
WorkingDirectory=/home/autosub/autosub/backend
EnvironmentFile=/home/autosub/autosub/backend/.env
# Concurrency=1 is intentional on the 2 vCPU box — keeps FFmpeg + the
# transcription HTTP call from saturating both cores at once. Bump if/when
# the host is upgraded.
ExecStart=/home/autosub/autosub/backend/.venv/bin/celery \
    -A app.tasks worker -c 1 --loglevel=info
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
